# coleting_client/etc/cron.d/coleting
30 0 * * * root /usr/lib/cliente/CopiasConTar.sh
