# coleting_client/etc/cron.d/coleting
30 0 * * * root /usr/lib/commserve/CopiasConTar.sh
0 1 * * * root /usr/lib/ommserve/ImportarLogs.sh
