# coleting_commserve/etc/systemd/system/grafana-deferred.timer
[Unit]
Description=Run deferred Grafana installer after boot

[Timer]
OnBootSec=1min
Unit=grafana-deferred.service

[Install]
WantedBy=timers.target
