# coleting_mediaagent/etc/cron.d/coleting
30 1 * * * root /usr/lib/mediaagent/ImportarBackups.sh
30 2 * * * root /usr/lib/mediaagent/Retencion.sh
